# scripts/sync/plan-fact-sync.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.foryou.plan-fact-sync</string>
  <key>WorkingDirectory</key>
  <string>/Users/vytvytskyi/foryou-dashboards</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/vytvytskyi/foryou-dashboards &amp;&amp; /usr/local/bin/npm run sync:plan-fact:bq</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/Users/vytvytskyi/foryou-dashboards/data/logs/plan-fact-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vytvytskyi/foryou-dashboards/data/logs/plan-fact-sync.error.log</string>
</dict>
</plist>